# services/rebuild-web.service


[Unit]
Description=Rebuild Web Assets
