# services/watch-file-changes.service

[Unit]
Description=Executes when /var/raw is updated

[Service]
Type=simple
ExecStart=/opt/sites/
