# services/admin/daily-tasks.service

[Unit]
Description=Dummy service for daily tasks

[Service]
Type=oneshot
ExecStart=echo "RUNNING DAILY TASKS"
